# system_info_service
#!/sbin/openrc-run

description="system info once after login"

start_pre() {
    /usr/sbin/system_info.sh
}

start() {
    return 0
}

stop() {
    return 0
}
